# docker/crontab
* * * * * /usr/local/bin/php /app/artisan schedule:run >> /app/storage/logs/schedule.log 2>&1
