# docker/crontab
* * * * * /usr/local/bin/php /app/artisan schedule:run >> /var/log/cron.log 2>&1
